# system_files/etc/systemd/system/htpc-cec.service
[Unit]
Description=CEC Daemon

[Service]
ExecStart=htpc-cecd

[Install]
WantedBy=multi-user.target
