# system_files/etc/systemd/system/htpc-cec.service
[Unit]
Description=CEC Daemon

[Service]
ExecStart=python3 -m htpc-cec-daemon

[Install]
WantedBy=multi-user.target
