# python-cec/cec-daemon.service
[Unit]
Description=CEC Daemon

[Service]
ExecStart=/home/bazzite/python-cec/dist/main/main

[Install]
WantedBy=multi-user.target